# appliance/root/etc/systemd/system/incus-spawn-diag.service
[Unit]
Description=Boot diagnostics dump

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/incus-spawn-diag
