# appliance/root/etc/systemd/system/incus-spawn-smoke-test.service
[Unit]
Description=Incus appliance smoke test
After=incus-spawn-vm.service
Requires=incus-spawn-vm.service
ConditionKernelCommandLine=isx.smoke_test

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/incus-spawn-smoke-test
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
