# appliance/root/etc/systemd/system/incus-spawn-diag.timer
[Unit]
Description=Boot diagnostics timer

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target
